# contrib/vm/compute_node/etc/cron.d/sagecell-cleantmp
*/30 * * * * root /usr/sbin/tmpreaper 2h /tmp/sagecell
