# contrib/vm/compute_node/etc/cron.d/sagecell-monitor
*/2 * * * * root /root/healthcheck.py http://localhost:8888 >> /root/healthcheck-cron.log 2>&1
