# contrib/vm/compute_node/etc/cron.d/sagecell-cleantmp
*/30 * * * * root /usr/sbin/tmpreaper --mtime-dir 2h /tmp
