# contrib/vm/upstart/sagecell-restart.conf
start on stopped sagecell-monitor

exec restart sagecell
